# reset/vpna_reset.conf
description	"Sabai VPNA reset listener"
author 		"Kirk Fort <kirk@sabaitechnology.com>"

start on runlevel [2345]


respawn

exec /usr/bin/perl /home/sabai/vpna_reset.pl
