# debian/etc/init/vpna_reboot.conf
description	"Sabai VPNA tunnel state loader"
author 		"Kirk Fort <kirk@sabaitechnology.com>"

start on filesystem and net-device-up IFACE=eth0

task

script
	sh /var/www/stat/vpn.command
end script
